# etc/supervisord.conf
[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:7101        ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile=./logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=./tmp/pids/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:7101 ; use an http:// url to specify an inet socket

[program:survey_application]
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s.pid
    start_worker
    --worker-class=go.apps.surveys.vumi_app.SurveyApplication
    --config=./config/poll.yaml
    --vhost=/go
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:haproxy]
command=/usr/sbin/haproxy -c ./config/haproxy.cfg -p ./tmp/pids/haproxy.pid
stdout_logfile=./logs/haproxy.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/haproxy.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autostart=true
autorestart=true

[program:go]
numprocs=4
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=django-admin.py run_gunicorn 127.0.0.1:807%(process_num)s
environment=DJANGO_SETTINGS_MODULE=go.settings,
            PYTHONPATH=.
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:celery]
command=django-admin.py celeryd --loglevel=INFO
environment=DJANGO_SETTINGS_MODULE=go.settings,
            PYTHONPATH=.
stdout_logfile=./logs/celeryd.log
stderr_logfile=./logs/celeryd.err
autostart=true
autorestart=true
startsecs=10

[program:command_dispatcher]
numprocs=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s.pid
    start_worker
    --worker-class=go.vumitools.api_worker.CommandDispatcher
    --config=./config/command_dispatcher.yaml
    --vhost=/go
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:application_message_dispatcher]
numprocs=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s.pid
    start_worker
    --worker-class=vumi.dispatchers.base.BaseDispatchWorker
    --config=./config/application_message_dispatcher.yaml
    --vhost=/go
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:bulk_message_application]
numprocs=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s.pid
    start_worker
    --worker-class=go.apps.bulk_message.vumi_app.BulkMessageApplication
    --config=./config/bulk_message_application.yaml
    --vhost=/go
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:smpp_transport]
numprocs=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s.pid
    start_worker
    --worker-class=vumi.transports.smpp.SmppTransport
    --config=./config/smpp_transport.yaml
    --vhost=/go
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:vumigo_gtalk_transports]
numprocs=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s.pid
    start_worker
    --worker-class=vumi.multiworker.MultiWorker
    --config=./config/gtalk_transports.yaml
    --vhost=/go
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:vumigo_router]
numprocs=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s.pid
    start_worker
    --worker-class=vumi.dispatchers.base.BaseDispatchWorker
    --config=./config/vumigo_router.yaml
    --vhost=/go
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true
