# setup_env/go_wikipedia_ussd_application.conf
# This file has been automatically generated by: 
# /Users/sdehaan/Documents/Repositories/go/go/base/management/commands/go_setup_env.py.

[program:wikipedia_ussd_application]
command = twistd -n --pidfile=./tmp/pids/wikipedia_ussd_application.pid start_worker --worker-class=go.apps.wikipedia.vumi_app.WikipediaApplication --config=setup_env/go_wikipedia_ussd_application.yaml
stdout_logfile = ./logs/%(program_name)s_%(process_num)s.log
stderr_logfile = ./logs/%(program_name)s_%(process_num)s.log

